# systemd/fetch-k8s-cert.timer
[Unit]
Description=Fetch K8S TLS certificate every 20 minutes

[Timer]
OnCalendar=*:0/20
AccuracySec=1s

[Install]
WantedBy=timers.target
